# docker/application/services.d/nginx/run
#!/usr/bin/with-contenv sh

if [ "$DISABLE_NGINX" != "1" ]
then
    exec nginx
fi
